# oasis-mod-tool/files/etc/uci-defaults/81_oasis-mod-tool
#!/bin/sh

changed=0

if [ "$(uci -q get oasis.support.local_tool)" != "1" ]; then
    uci set oasis.support.local_tool='1'
    changed=1
fi

if [ "$(uci -q get oasis.support.remote_mcp_server)" != "1" ]; then
    uci set oasis.support.remote_mcp_server='1'
    changed=1
fi

if [ "$changed" -eq 1 ]; then
    uci commit oasis
fi

if [ -x /usr/bin/oasis ]; then
    for manifest in /etc/oasis/tool-manifest.d/*.json; do
        [ -f "$manifest" ] || continue
        /usr/bin/oasis manifest apply --yes "$manifest" >/dev/null 2>&1 || true
    done
fi

exit 0
