# oasis/files/etc/init.d/oasis.init
#!/bin/sh /etc/rc.common

START=99

start() {

    mkdir -p /tmp/oasis/apply
    # rollback timer start!!
    lua /usr/bin/oasis_rollback boot &
}
