# oasis-mod-spring/files/etc/init.d/spring.init
#!/bin/sh /etc/rc.common

START=99

bootstart=$(uci -q get spring.ctrl.bootstart)

start() {

	mkdir -p /tmp/spring

	if [ "${bootstart}" = "1" ];then
        /usr/bin/springd
		echo "springd start!"
	fi
}

stop() {
	killall /usr/bin/springd
	echo 'springd stop'
}

restart() {
	stop > /dev/null
	sleep 1
	start > /dev/null
	echo 'springd restart'
}